# network-report.service
[Unit]
Description=Network Service Report
After=network.target avahi-daemon.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/network-report.sh /var/www/html/reports/network-report.html
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target 